# lib/schema-registry/etc/init/schema-registry.conf
# sudo add-apt-repository "deb [arch=amd64] https://packages.confluent.io/deb/3.2 stable main"
# apt-get update && apt-get install confluent-schema-registry

description "Confluent Schema Registry"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
post-stop exec sleep 15

umask 007

kill timeout 300
setuid root
setgid root

#env SCHEMA_REGISTRY_OPTS=-Djava.security.auth.login.config=/etc/schema-registry/jaas.conf
#export SCHEMA_REGISTRY_OPTS
exec /usr/bin/schema-registry-start /etc/schema-registry/schema-registry.properties