# env/example-app/etc/init/example-app.conf
description "Exmple Scala App"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 2 5

umask 007

kill timeout 300
setuid root
setgid root

exec /usr/bin/java -jar -Dlogback.configurationFile=/etc/example-app/logback.xml \
  /opt/example-api/build/libs/example-app-1.1.2-all.jar \
  /etc/example-app/akka.conf