# lib/zookeeper/etc/init/zookeeper.conf
description "ZooKeeper"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
post-stop exec sleep 15

umask 007

kill timeout 300
setuid root
setgid root

exec /opt/kafka/current/bin/zookeeper-server-start.sh /etc/zookeeper/zoo.cfg
